# .config/autostart/fastfetch_wall.desktop
[Desktop Entry]
Type=Application
Name=Conky
Exec=sh -c "sleep 5 && ~/.config/conky/autostart.sh"
Comment=Inicia o Conky
Terminal=false